# ollama-distributed/deploy/systemd/ollamacron.service
[Unit]
Description=Ollamacron - Distributed Ollama Platform
Documentation=https://github.com/ollama/ollama-distributed
After=network.target
Wants=network.target

[Service]
Type=simple
User=ollama
Group=ollama
ExecStart=/usr/local/bin/ollamacron node \
  --config /etc/ollamacron/config.yaml \
  --log-level info \
  --log-format json
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/ollamacron /var/log/ollamacron
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
LimitAS=infinity
LimitCORE=0

# Working directory
WorkingDirectory=/var/lib/ollamacron

# Environment variables
Environment=OLLAMACRON_CONFIG=/etc/ollamacron/config.yaml
Environment=OLLAMACRON_DATA_DIR=/var/lib/ollamacron
Environment=OLLAMACRON_LOG_LEVEL=info
Environment=OLLAMACRON_LOG_FORMAT=json

[Install]
WantedBy=multi-user.target